# root/etc/s6/gitea/run
#!/bin/bash
[[ -f ./setup ]] && source ./setup

pushd ${GITEA_HOME}/gitea > /dev/null
    exec su-exec ${USER} ${GITEA_HOME}/gitea/gitea web
popd
